# etc/service/rsyslog/run
#!/bin/bash -e

exec /usr/sbin/rsyslogd -n -c5 2>&1
