# etc/service/wwsympa/run
#!/bin/bash -e

# Sleep so that main Sympa process can start and prepare everything.
sleep 10

export FCGI_SOCKET_PATH="127.0.0.1:9000"
echo fastcgi_pass "127.0.0.1:9000;" > /etc/nginx/fastcgi_host

exec chpst -u sympa:sympa /usr/lib/cgi-bin/sympa/wwsympa.fcgi 2>&1
