# etc/service/sympa_task_manager/run
#!/bin/bash -e

# Sleep so that main Sympa process can start and prepare everything.
sleep 10

exec chpst -u sympa:sympa /usr/lib/sympa/bin/task_manager.pl --foreground
